# docker/s6-rc.d/acme-renew/run
#!/command/with-contenv sh
ACME="${HOME}/.acme.sh/acme.sh"
ACME_HOME="${HOME}/.acme.sh"
PORT_MIN=49152
PORT_MAX=65535

is_port_free() {
  ! nc -z 127.0.0.1 "$1" 2>/dev/null
}

configured_port() {
  case "${ACME_HTTP_PORT:-}" in
    ''|*[!0-9]*)
      echo 51823
      ;;
    *)
      if [ "${ACME_HTTP_PORT}" -ge "$PORT_MIN" ] && [ "${ACME_HTTP_PORT}" -le "$PORT_MAX" ]; then
        echo "${ACME_HTTP_PORT}"
      else
        echo 51823
      fi
      ;;
  esac
}

pick_free_port() {
  preferred=$(configured_port)
  if is_port_free "$preferred"; then
    echo "$preferred"
    return
  fi
  port=$PORT_MIN
  while [ "$port" -le "$PORT_MAX" ]; do
    if [ "$port" -ne "$preferred" ] && is_port_free "$port"; then
      echo "$port"
      return
    fi
    port=$((port + 1))
  done
  echo "$preferred"
}

patch_acme_port() {
  free_port=$(pick_free_port)
  for conf in "$ACME_HOME"/*/*.conf; do
    [ -f "$conf" ] || continue
    if grep -q "Le_Httpport=" "$conf" 2>/dev/null; then
      sed -i "s/Le_Httpport='[^']*'/Le_Httpport='$free_port'/" "$conf"
    fi
  done
}

while true; do
  if [ -f "$ACME" ]; then
    patch_acme_port
    "$ACME" --cron --home "$ACME_HOME"
  fi
  sleep 86400
done
